# scripts_installer/lin/deb/pulse-agent-linux/debian/pulse-agent-linux.postinst.in
#!/bin/bash
VNC_PORT="@@VNC_PORT@@"
SSH_PORT="@@SSH_PORT@@"

### Pulse account

echo
echo "#########################"
echo "1. CREATING PULSE ACCOUNT"
echo "#########################"
echo

# Create medulla user
adduser --system --group --home /home/medulla --shell /bin/rbash --disabled-password medulla
if ! [ -d "/home/medulla/.ssh" ]; then
    echo "Create missing SSH profile ..."
    mkdir -p /home/medulla/.ssh
fi



### SSH Key

echo
echo "############"
echo "2. SETUP SSH"
echo "############"
echo

# Copy SSH keys
chown -R medulla: /home/medulla
chmod -R 700 /home/medulla
chmod 600 /home/medulla/.ssh/authorized_keys

sed -i "s/^#?Port .*$/Port ${SSH_PORT}/" /etc/ssh/sshd_config

# Open firewall port 22 - TODO
# if shorewall
# echo "ACCEPT net fw tcp 22" >> /etc/shorewall/rules
# service shorewall restart
#
# if ufw
# ufw enable
# sudo ufw allow 5900



echo
echo "###########################"
echo "3. SETUP SUDO for medulla"
echo "###########################"
echo
echo "medulla ALL=(ALL) NOPASSWD: /usr/bin/rsync" > /etc/sudoers.d/medulla



### Setup rdp server

echo
echo "###################"
echo "4. SETUP RDP SERVER"
echo "###################"
echo

systemctl enable xrdp

# Set permissions on the ssl keys
usermod -aG ssl-cert xrdp
chown root:ssl-cert /etc/xrdp/*.pem

# Start the service
systemctl restart xrdp

# Open firewall port 3389 - TODO
# if shorewall
# echo "ACCEPT net fw tcp 3389" >> /etc/shorewall/rules
# service shorewall restart
#
# if ufw
# ufw enable
# sudo ufw allow 3389
#

echo
echo "##########################"
echo "5. Install the certificate"
echo "##########################"
echo

update-ca-certificates

