# scripts_installer/lin/deb/pulse-agent-linux/debian/pulse-agent-linux.postinst.in
#!/bin/bash
VNC_PORT="@@VNC_PORT@@"
SSH_PORT="@@SSH_PORT@@"

### Pulse account

echo
echo "#########################"
echo "1. CREATING PULSE ACCOUNT"
echo "#########################"
echo

# Create pulse user
adduser --system --group --home /var/lib/pulse2 --shell /bin/rbash --disabled-password pulseuser
if ! [ -d "/var/lib/pulse2/.ssh" ]; then
    echo "Create missing SSH profile ..."
    mkdir -p /var/lib/pulse2/.ssh
fi



### SSH Key

echo
echo "############"
echo "2. SETUP SSH"
echo "############"
echo

# Copy SSH keys
chown -R pulseuser: /var/lib/pulse2
chmod -R 700 /var/lib/pulse2
chmod 600 /var/lib/pulse2/.ssh/authorized_keys

sed -i "s/^#?Port .*$/Port ${SSH_PORT}/" /etc/ssh/sshd_config

# Open firewall port 22 - TODO
# if shorewall
# echo "ACCEPT net fw tcp 22" >> /etc/shorewall/rules
# service shorewall restart
#
# if ufw
# ufw enable
# sudo ufw allow 5900



echo
echo "###########################"
echo "3. SETUP SUDO for pulseuser"
echo "###########################"
echo
echo "pulseuser ALL=(ALL) NOPASSWD: /usr/bin/rsync" > /etc/sudoers.d/pulseuser



### Setup rdp server

echo
echo "###################"
echo "4. SETUP RDP SERVER"
echo "###################"
echo

systemctl enable xrdp

# Open firewall port 3389 - TODO
# if shorewall
# echo "ACCEPT net fw tcp 3389" >> /etc/shorewall/rules
# service shorewall restart
#
# if ufw
# ufw enable
# sudo ufw allow 3389
