# debian/pulse-xmpp-master-substitute.postinst
#!/bin/sh
set -e

systemctl daemon-reload

case "$1" in
    configure)
        if systemctl is-enabled pulse-xmpp-master-substitute-inventory &>/dev/null ; then
            systemctl restart pulse-xmpp-master-substitute-inventory
        fi
        if systemctl is-enabled pulse-xmpp-master-substitute-registration &>/dev/null ; then
            systemctl restart pulse-xmpp-master-substitute-registration
        fi
        if systemctl is-enabled pulse-xmpp-master-substitute-assessor &>/dev/null ; then
            systemctl restart pulse-xmpp-master-substitute-assessor
        fi
        if systemctl is-enabled pulse-xmpp-master-substitute-deployment &>/dev/null ; then
            systemctl restart pulse-xmpp-master-substitute-deployment
        fi
        if systemctl is-enabled pulse-xmpp-master-substitute-subscription &>/dev/null ; then
            systemctl restart pulse-xmpp-master-substitute-subscription
        fi
        if systemctl is-enabled pulse-xmpp-master-substitute-logger &>/dev/null ; then
            systemctl restart pulse-xmpp-master-substitute-logger
        fi
        if systemctl is-enabled pulse-xmpp-master-substitute-all &>/dev/null ; then
            systemctl restart pulse-xmpp-master-substitute-all
        fi
        ;;
esac

#DEBHELPER#

exit 0
