# pulse_xmpp_agent/config/systemd/pulse-xmpp-agent-log.service
[Unit]
Description=Pulse2 XMPP Agent
After=samba.service slapd.service ejabberd.service
Requires=ejabberd.service

[Service]
#PIDFile=/var/run/pulse-xmpp-agent-log.pid
#User=pulse
#Group=pulse
Type=forking
ExecStart=/usr/lib/python3/site-packages/pulse_xmpp_agent/pulse-xmpp-agent-log.py -d
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
